# teleport-updater.service
[Unit]
Description=Teleport Updater

[Service]
Type=oneshot
ExecStart=/usr/local/bin/teleport-updater
